# service/sshd/run
#!/bin/sh

exec /usr/sbin/sshd -D 2>&1

